# service/pisphere.timer.template
[Unit]
Description=PiSphere Image Capture Timer

[Timer]
# Trigger run.sh every INTERVAL_MINUTES within the capture window
OnCalendar=*-*-* START_TIME/INTERVAL_MINUTES
Persistent=true

[Install]
WantedBy=timers.target
